# deploy/ptspgt2/systemd/multica-daemon.service
[Unit]
Description=multica agent daemon
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root
Environment=HOME=/root
Environment=IS_SANDBOX=1
Environment=PATH=/root/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/snap/bin
ExecStart=/usr/local/bin/multica daemon start --foreground
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
